# rpi5/rpi5-install/etc/udev/rules.d/80-ethernet-autodetect.rules
# Trigger MANET ethernet role detection on carrier changes.
SUBSYSTEM=="net", KERNEL=="end0", ACTION=="change", RUN+="/bin/systemctl start ethernet-autodetect.service"
