# rpi5/rpi5-install/etc/udev/rules.d/80-usb-ethernet.rules
# Trigger USB ethernet watch when a USB-backed network interface appears or disappears
ACTION=="add",    SUBSYSTEM=="net", ENV{ID_BUS}=="usb", RUN+="/usr/local/bin/usb-ethernet-watch.sh add $name"
ACTION=="remove", SUBSYSTEM=="net", ENV{ID_BUS}=="usb", RUN+="/usr/local/bin/usb-ethernet-watch.sh remove $name"
