# rpi5/rpi5-live/78f7/etc/systemd/system/mesh-default-route-fix.service
[Unit]
Description=Set mesh default route after boot
After=network-online.target batman-enslave.service node-manager.service gateway-route-manager.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mesh-default-route-fix.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
