# rpi5/rpi5-tools/etc/systemd/system/mesh-clone-identity.service
[Unit]
Description=Reset cloned mesh node identity when hardware MAC changes
After=local-fs.target
Before=batman-enslave.service alfred.service node-manager.service gateway-route-manager.service syncthing@radio.service hostapd.service dnsmasq.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mesh-clone-identity.sh

[Install]
WantedBy=multi-user.target
