# rpi5/rpi5-tools/etc/systemd/system/ssh-recovery.service
[Unit]
Description=Ensure MANET SSH access is available
After=local-fs.target
Before=network-online.target batman-enslave.service alfred.service node-manager.service hostapd.service dnsmasq.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-recovery.sh

[Install]
WantedBy=multi-user.target
