# rpi5/rpi5-install/etc/systemd/system/ethernet-autodetect.service
[Unit]
Description=MANET Ethernet Hotplug Auto Detection
After=systemd-networkd.service batman-enslave.service
Wants=systemd-networkd.service
ConditionPathExists=/usr/local/bin/manet-uplink-dispatch.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/manet-uplink-dispatch.sh reconcile
TimeoutStartSec=45

[Install]
WantedBy=multi-user.target
