# rpi5/rpi5-install/etc/systemd/system/mesh-mdns-update.timer
[Unit]
Description=Refresh local mDNS aliases periodically

[Timer]
OnBootSec=45
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
